# supervisord/CONTROL/postinst
#!/bin/sh
pip install supervisor==3.0

systemctl enable supervisord
systemctl start supervisord
