# supervisord/CONTROL/postinst
#!/bin/sh
pip install supervisor==3.0
update-rc.d supervisord defaults
/etc/init.d/supervisord start
