# supervisord/CONTROL/prerm
#!/bin/sh
systemctl stop supervisord
systemctl disable supervisord
pip uninstall -y supervisor
